# mybase/reset.cron
00 13 * * * ~/github/mybase/mybase/resetTimeStamp.sh
